# init/alpine/gateway-owm
#!/sbin/openrc-run
command=/usr/local/bin/gateway-owm
command_args="-log-file=/var/log/eria/${RC_SVCNAME}.log"

conf=/usr/local/etc/eria

pidfile=/var/run/${RC_SVCNAME}.pid
name="ERIA OWM Gateway Daemon"

description="ERIA Project, gateway for OWM service"

start() {
  ebegin "Starting ${name}"
  start-stop-daemon --background --start --exec \
    env ERIA_CONF_PATH=${conf} \
    ${command} \
    --make-pidfile --pidfile ${pidfile} \
    -- ${command_args} \
  eend $?
}

stop() {
  ebegin "Stopping ${name}"
  start-stop-daemon --stop \
    --exec ${command} \
    --pidfile ${pidfile} \
  eend $?
}
